# com.prorabbot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.prorabbot</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPLACE_WITH_FULL_PATH/start_bot.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_FULL_PATH</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>REPLACE_WITH_FULL_PATH/logs/bot.log</string>

  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_FULL_PATH/logs/bot_error.log</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
